# conf.d/jdownloader2.conf
[program:jdownloader2]
command=java -jar JDownloader.jar
autorestart=true
