# conf.d/jdownloader2.conf
[program:jdownloader2]
command=bash -c "java -jar JDownloader.jar; sleep 10"
autorestart=true
