# app/conf.d/jdownloader2.conf
[program:jdownloader2]
command=bash /start-jdownloader.sh
autorestart=true
